# debian/galette.cron.daily
#!/bin/sh

# Exit if pre-requesites are not met
if [ ! -e /usr/share/galette/reminder.php ] || ! which php5 >/dev/null 2>&1; then
    exit 0
fi

# Send the reminder mails for members who have to renew their membership
su - www-data -s /bin/sh -c "php5 -d error_reporting=22519 /usr/share/galette/reminder.php"
