# debian/galette.cron.daily
#!/bin/sh

# Send the reminder mails for members who have to renew their membership
su - www-data -s /bin/sh -c "php5 /usr/share/galette/reminder.php"
